# roles/news/files/news.cron
SHELL=/bin/bash
PATH=/sbin:/bin/:/usr/sbin:/usr/bin
*/5 * * * * root /usr/local/bin/news pull >/dev/null