# roles/efingerd/files/finger@.service
[Unit]
Description=Finger Per-Connection Server

[Service]
User=nobody
ExecStart=/usr/local/sbin/efingerd
StandardInput=socket
